# crontab_service/crontab
*/3 5-23 * * * curl "njt-bus-time-flask.railway.internal:8080/push"
