# crontab_service/crontab
* * * * * curl "njt-bus-time-flask.railway.internal/cron"
